# mailinglists/init.d/sympa2
#!/sbin/openrc-run
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
	need rsyslog
}

start() {
	/etc/init.d/sympa start
}

stop() {
	/etc/init.d/sympa stop
}

start_pre() {
	/home/sympa/bin/sympa.pl --health_check
	ebegin "Starting sympa"
}
